# n8n-workflows/deployment/n8n-failure-alert.service
[Unit]
Description=n8n Service Failure Email Alert
Documentation=https://github.com/yourusername/spa-kiosk

[Service]
Type=oneshot
# Set the alert email address (override with systemctl edit n8n-failure-alert.service)
Environment=ALERT_EMAIL=admin@example.com
ExecStart=/usr/local/bin/n8n-failure-alert.sh

# Run as root to access systemctl status
User=root

# Timeout after 30 seconds
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
